# .config/systemd/user/dotfiles-sync.timer
[Unit]
Description=Dotfiles auto sync timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
